# btrfs-home-snapshot/btrfs-home-snapshot.cron
@reboot root /usr/local/scripts/btrfs-home-snapshot.sh